# roles/llama_cpp/templates/llama-swap.service.j2
{{ ansible_managed | comment }}
[Unit]
Description=llama-swap — OpenAI-compatible proxy over llama.cpp (ROCm light tier)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ llama_cpp_user }}
Group={{ llama_cpp_group }}
WorkingDirectory={{ llama_cpp_data_dir }}
# The ROCm release binary and its bundled .so files share the install dir.
Environment=LD_LIBRARY_PATH={{ llama_cpp_install_dir }}
{% if llama_cpp_gpu | bool %}
Environment=HSA_OVERRIDE_GFX_VERSION={{ llama_cpp_hsa_override_gfx_version }}
{% endif %}
ExecStart={{ llama_cpp_swap_bin }} --config {{ llama_cpp_config_file }} --listen 0.0.0.0:{{ llama_cpp_api_port }}
# Restart policy is managed by the systemd_restart_policy role (99-restart-policy.conf).

[Install]
WantedBy=multi-user.target
