# roles/postgres/templates/pg-backup.service.j2
[Unit]
# {{ ansible_managed }}
Description=PostgreSQL pg_dump backup to the persistent mount
After=postgresql@{{ postgres_cluster_version }}-{{ postgres_cluster_name }}.service
Wants=postgresql@{{ postgres_cluster_version }}-{{ postgres_cluster_name }}.service

[Service]
Type=oneshot
User=postgres
Group=postgres
ExecStart=/usr/local/bin/pg-backup.sh
# Backups are local disk I/O; never let a slow dump wedge the timer forever.
TimeoutStartSec=3600
