# roles/postgres/templates/pgbackrest-backup@.service.j2
[Unit]
# {{ ansible_managed }}
# Instanced unit: %i is the backup type (full/incr), so one service definition
# serves every schedule. Runs as postgres over the local socket (peer auth).
Description=pgBackRest %i backup for stanza {{ postgres_pgbackrest_stanza }}
After=postgresql@{{ postgres_cluster_version }}-{{ postgres_cluster_name }}.service
Wants=postgresql@{{ postgres_cluster_version }}-{{ postgres_cluster_name }}.service

[Service]
Type=oneshot
User=postgres
Group=postgres
ExecStart=/usr/bin/pgbackrest --stanza={{ postgres_pgbackrest_stanza }} --type=%i backup
# A base backup to remote S3 can take a long time; never wedge it on a timeout.
TimeoutStartSec=0
