# roles/agent_sandbox/templates/spool-prune.timer.j2
# {{ ansible_managed }}
[Unit]
Description=Daily prune of the agent-sandbox transcript spool

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
