# roles/agent_sandbox/templates/spool-prune.service.j2
# {{ ansible_managed }}
[Unit]
Description=Prune agent-sandbox transcript spool runs older than {{ agent_sandbox_spool_retention_days }} days

[Service]
Type=oneshot
# Each <run-id> dir is a self-contained run; drop whole runs past the window.
ExecStart=/usr/bin/find {{ agent_sandbox_spool_dir }} -mindepth 1 -maxdepth 1 -type d -mtime +{{ agent_sandbox_spool_retention_days }} -exec rm -rf {} +
