# roles/llamaindex/templates/llamaindex-index.service.j2
[Unit]
Description=Rebuild the homelab RAG index in Qdrant
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Router key and Qdrant key are baked into config.yaml at converge time. Only
# the optional Cloudflare Access service token for the private docs source is
# read from the environment at runtime; point EnvironmentFile at wherever it
# lives to enable that source (otherwise it is skipped, non-fatally).
{% if llamaindex_index_environment_file is defined %}
EnvironmentFile={{ llamaindex_index_environment_file }}
{% endif %}
ExecStart={{ llamaindex_data_dir }}/venv/bin/python {{ llamaindex_data_dir }}/index_docs.py --config {{ llamaindex_data_dir }}/config.yaml
