# roles/agent_exec/templates/agent-exec.service.j2
{{ ansible_managed | comment }}
[Unit]
Description=agent-exec — CrewAI / LangChain agent runtime (OpenLLMetry traced)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ agent_exec_user }}
Group={{ agent_exec_user }}
WorkingDirectory={{ agent_exec_home }}
Environment=OTEL_EXPORTER_OTLP_ENDPOINT={{ ai_orchestration_otel_endpoint }}
Environment=OTEL_EXPORTER_OTLP_PROTOCOL=http/protobuf
Environment=OTEL_SERVICE_NAME={{ agent_exec_service }}
# Two-tier model fabric (re-rendered each converge, so the backend is swappable):
# both tiers now reach the LiteLLM router, which routes by model alias.
Environment=AGENT_EXEC_MODEL_LARGE_BASE_URL={{ ai_orchestration_model_large_base_url }}
Environment=AGENT_EXEC_MODEL_LIGHT_BASE_URL={{ ai_orchestration_ollama_base_url }}
Environment=AGENT_EXEC_MODEL_LARGE={{ agent_exec_model_large }}
Environment=AGENT_EXEC_MODEL_LIGHT={{ agent_exec_model_light }}
Environment=AGENT_EXEC_INTERVAL_SECONDS={{ agent_exec_interval_seconds }}
Environment=OPENAI_API_KEY={{ agent_exec_api_key }}
ExecStart={{ agent_exec_venv }}/bin/python {{ agent_exec_home }}/main.py
# Inline restart policy: this role must not edit the shared systemd_restart_policy
# role, so it sets its own restart-on-failure here.
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
