# roles/agent_exec/templates/agent-exec.service.j2
{{ ansible_managed | comment }}
[Unit]
Description=agent-exec — CrewAI / LangChain agent runtime (OpenLLMetry traced)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ agent_exec_user }}
Group={{ agent_exec_user }}
WorkingDirectory={{ agent_exec_home }}
Environment=OTEL_EXPORTER_OTLP_ENDPOINT={{ ai_orchestration_otel_endpoint }}
Environment=OTEL_EXPORTER_OTLP_PROTOCOL=http/protobuf
Environment=OTEL_SERVICE_NAME={{ agent_exec_service }}
ExecStart={{ agent_exec_venv }}/bin/python {{ agent_exec_home }}/main.py
# Inline restart policy: this role must not edit the shared systemd_restart_policy
# role, so it sets its own restart-on-failure here.
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
