# roles/download_vpn/templates/download-vpn-validate.service.j2
[Unit]
Description=download-vpn fail-closed killswitch validator (one-shot, timer-driven)
After=download-vpn-wg.service

[Service]
Type=oneshot
# Needs the qBittorrent admin password to query the WebUI API.
EnvironmentFile={{ download_vpn_natpmp_env_path }}
ExecStart={{ download_vpn_validator_script_path }}
