# roles/download_vpn/templates/download-vpn-natpmp.service.j2
[Unit]
Description=download-vpn Proton NAT-PMP port-forward keeper
After=download-vpn-wg.service qbittorrent-nox.service
Requires=download-vpn-killswitch.service download-vpn-wg.service
BindsTo=download-vpn-wg.service

[Service]
Type=simple
User={{ download_vpn_user }}
Group={{ download_vpn_group }}
# QBITTORRENT_ADMIN_PASSWORD is delivered via a root-only EnvironmentFile so
# the script can authenticate to the WebUI API.
EnvironmentFile={{ download_vpn_natpmp_env_path }}
ExecStart={{ download_vpn_natpmp_script_path }}
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
