# roles/download_vpn/templates/download-vpn-lanroute.service.j2
[Unit]
Description=download-vpn LAN-reply policy routing (cross-subnet web UI reachability)
Documentation=https://github.com/JacobPEvans/ansible-proxmox-apps
# Order after wg0 so wg-quick's own fwmark policy rules already exist and we
# layer a higher-priority rule above them; network-online ensures the LAN NIC is
# configured. Deliberately NOT Requires=download-vpn-wg.service: this reply
# routing is independent of the tunnel, so the web UIs stay reachable for
# debugging even during a VPN outage. The killswitch (not this unit) is the
# egress security boundary, so keeping it up never enables a leak.
After=network-online.target download-vpn-wg.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes

# 0. Wait for the LAN gateway to become routable VIA THE LAN NIC before adding
#    policy routes. In an unprivileged LXC, network-online.target can fire BEFORE
#    the DHCP lease installs eth0's on-link route, so the route adds below would
#    fail with "Nexthop has invalid gateway" — and since this is a oneshot, it
#    would stay failed (table empty, web UIs unreachable cross-subnet until a
#    manual restart). The probe is pinned to `dev {{ download_vpn_lan_interface }}`
#    on purpose: this unit is ordered After=download-vpn-wg.service, so wg0's
#    default route is already up, and an unscoped `ip route get` would resolve the
#    gateway via wg0 and pass prematurely while eth0 is still lease-less. Bounded
#    to ~30s so a genuine network outage fails the unit cleanly instead of hanging
#    boot. Pure readiness gate — opens no egress path.
ExecStartPre=/usr/bin/env sh -c 'n=0; until ip route get {{ download_vpn_lan_gateway }} dev {{ download_vpn_lan_interface }} >/dev/null 2>&1; do n=$((n+1)); [ "$n" -ge 30 ] && exit 1; sleep 1; done'

# 1. (Re)load the connmark mangle table from a clean slate. Delete first (ignore
#    if absent) so a restart never hits "chain already exists". This must succeed
#    on load — without the marking the policy route below is inert.
ExecStartPre=-/usr/sbin/nft delete table inet vpn_lanroute
ExecStartPre=/usr/sbin/nft -f {{ download_vpn_mangle_ruleset_path }}

# 2. Policy route: reply packets carrying our fwmark resolve out the LAN NIC via
#    a dedicated table; everything unmarked keeps using wg0 (or is dropped by the
#    killswitch). Delete-then-add the rule for idempotency; `route replace` is
#    inherently idempotent. `env ip` avoids hardcoding the iproute2 path.
#
#    The table routes ONLY RFC1918 ranges (legitimate marked replies only ever
#    target LAN clients, never public IPs) and ends in `unreachable default`. So a
#    packet that somehow carries the mark but is destined for the public internet
#    finds no route here and fails CLOSED — clearnet egress via the LAN table is
#    structurally impossible, not merely policy-gated by the killswitch.
ExecStartPre=-/usr/bin/env ip rule del fwmark {{ download_vpn_lanroute_mark }} table {{ download_vpn_lanroute_table }} priority 1000
ExecStart=/usr/bin/env ip rule add fwmark {{ download_vpn_lanroute_mark }} table {{ download_vpn_lanroute_table }} priority 1000
{% for net in download_vpn_lanroute_reply_networks %}
ExecStart=/usr/bin/env ip route replace {{ net }} via {{ download_vpn_lan_gateway }} dev {{ download_vpn_lan_interface }} table {{ download_vpn_lanroute_table }}
{% endfor %}
ExecStart=/usr/bin/env ip route replace unreachable default table {{ download_vpn_lanroute_table }}

# Teardown: remove the rule, flush the table, drop the mangle marking — leaving
# no standing non-VPN routing state.
ExecStop=-/usr/bin/env ip rule del fwmark {{ download_vpn_lanroute_mark }} table {{ download_vpn_lanroute_table }} priority 1000
ExecStop=-/usr/bin/env ip route flush table {{ download_vpn_lanroute_table }}
ExecStop=-/usr/sbin/nft delete table inet vpn_lanroute

[Install]
WantedBy=multi-user.target
