# roles/download_vpn/templates/download-vpn-validate.timer.j2
[Unit]
Description=Run the download-vpn killswitch validator every {{ download_vpn_validator_interval }}
Documentation=https://github.com/JacobPEvans/ansible-proxmox-apps

[Timer]
# Start shortly after boot, then on a fixed cadence. Persistent=true catches
# missed runs (deadman semantics rely on regular execution).
OnBootSec=90s
OnUnitActiveSec={{ download_vpn_validator_interval }}
AccuracySec=15s
Persistent=true

[Install]
WantedBy=timers.target
