# roles/download_vpn/templates/download-vpn-wg.service.j2
[Unit]
Description=download-vpn WireGuard tunnel ({{ download_vpn_wg_interface }} -> Proton)
After=network-online.target download-vpn-killswitch.service
Wants=network-online.target
# Never bring the tunnel up without the killswitch already loaded.
Requires=download-vpn-killswitch.service

[Service]
Type=oneshot
RemainAfterExit=yes
# wg-quick honours `Table = auto` in wg0.conf: it installs the IPv4 default
# route via wg0 *and* an fwmark + policy-routing table so wg's encapsulated
# UDP to the Proton endpoint resolves via the original LAN default. Do not
# add a manual `ip route replace default dev wg0` here — it would clobber
# wg-quick's setup and recursively route the handshake back into the tunnel.
ExecStart=/usr/bin/wg-quick up {{ download_vpn_wg_interface }}
ExecStop=/usr/bin/wg-quick down {{ download_vpn_wg_interface }}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
