# roles/download_vpn/templates/flaresolverr.service.j2
[Unit]
Description=FlareSolverr CloudFlare challenge solver (VPN-locked)
Documentation=https://github.com/FlareSolverr/FlareSolverr
After=download-vpn-wg.service
# Boot-safe: FlareSolverr reaches CF-gated trackers via the VPN. Its clearance
# cookie is bound by CloudFlare to the exit IP, so it must never run without the
# killswitch AND the tunnel (a solve on a non-VPN IP would both leak and be
# useless to VPN-egressing Prowlarr).
Requires=download-vpn-killswitch.service download-vpn-wg.service
BindsTo=download-vpn-wg.service

[Service]
Type=simple
User={{ download_vpn_flaresolverr_user }}
Group={{ download_vpn_flaresolverr_user }}
# Loopback bind — only same-host Prowlarr consumes it; no LAN exposure.
Environment=HOST={{ download_vpn_flaresolverr_bind_host }}
Environment=PORT={{ download_vpn_flaresolverr_port }}
Environment=LOG_LEVEL=info
ExecStart={{ download_vpn_flaresolverr_install_dir }}/flaresolverr
Restart=on-failure
RestartSec=10
# Hard ceiling: headless Chromium is the guaranteed OOM victim, never the
# co-hosted qBittorrent (session state) or the Plex-serving host.
MemoryMax={{ download_vpn_flaresolverr_memory_max }}
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
