# roles/download_vpn/templates/prowlarr.service.j2
[Unit]
Description=Prowlarr indexer manager (VPN-locked)
Documentation=https://wiki.servarr.com/prowlarr
After=download-vpn-wg.service
# Boot-safe: Prowlarr reaches indexers via the VPN; never start without the
# killswitch AND the tunnel.
Requires=download-vpn-killswitch.service download-vpn-wg.service
BindsTo=download-vpn-wg.service

[Service]
Type=simple
User={{ download_vpn_user }}
Group={{ download_vpn_group }}
ExecStart={{ download_vpn_prowlarr_install_dir }}/Prowlarr -nobrowser -data={{ download_vpn_prowlarr_data_dir }}
Restart=on-failure
RestartSec=10
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
