# roles/download_vpn/templates/qbittorrent-nox.service.j2
[Unit]
Description=qBittorrent-nox (VPN-locked)
Documentation=https://github.com/qbittorrent/qBittorrent
After=download-vpn-wg.service
# Boot-safe: qBittorrent never starts without the killswitch AND the tunnel.
Requires=download-vpn-killswitch.service download-vpn-wg.service
# If the tunnel unit stops, take qBittorrent down with it.
BindsTo=download-vpn-wg.service

[Service]
Type=exec
User={{ download_vpn_user }}
Group={{ download_vpn_group }}
ExecStart=/usr/bin/qbittorrent-nox --webui-port={{ download_vpn_qbittorrent_web_port }}
# Resilience: a crashed daemon (e.g. disk-full abort) must come back on its
# own instead of staying dead until the next converge. 30s of backoff keeps a
# hard-failing daemon from hammering the disk it just filled.
Restart=on-failure
RestartSec=30s
# Hardening
NoNewPrivileges=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
