# roles/authelia/templates/authelia.service.j2
{{ ansible_managed | comment }}
[Unit]
Description=Authelia SSO portal / forwardAuth provider
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ authelia_user }}
Group={{ authelia_group }}
WorkingDirectory={{ authelia_install_dir }}
ExecStart={{ authelia_binary_path }} --config {{ authelia_config_file }}
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
