# roles/openbao/templates/openbao-snapshot.timer.j2
{{ ansible_managed | comment }}
# Fire the OpenBao snapshot every {{ openbao_snapshot_interval }}. Rendered by
# the openbao role.

[Unit]
Description=Take an OpenBao raft snapshot every {{ openbao_snapshot_interval }}
Documentation=https://github.com/JacobPEvans/ansible-proxmox-apps

[Timer]
# Start a few minutes after boot (let the node join + unseal first), then on a
# fixed cadence. Persistent=true runs a missed snapshot after downtime.
OnBootSec=300s
OnUnitActiveSec={{ openbao_snapshot_interval }}
AccuracySec=1m
Persistent=true

[Install]
WantedBy=timers.target
