# roles/llm_router/templates/litellm-rotate.timer.j2
{{ ansible_managed | comment }}
# Fire the {{ _rotation_phase }}-brain rotation at {{ _rotation_oncalendar }}
# (docs/BRAIN_ROTATION.md). Rendered by the llm_router role.

[Unit]
Description=Daily fabric-brain rotation to the {{ _rotation_phase }} model

[Timer]
OnCalendar={{ _rotation_oncalendar }}
AccuracySec=1s
# Run a rotation missed during downtime on the next boot, so the live phase is
# always correct for the current UTC half-day even after an outage.
Persistent=true

[Install]
WantedBy=timers.target
