# roles/llm_router/templates/litellm.service.j2
{{ ansible_managed | comment }}
[Unit]
Description=LiteLLM proxy — the LLM fabric front door (OpenAI-compatible)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ llm_router_user }}
Group={{ llm_router_group }}
WorkingDirectory={{ llm_router_home }}
EnvironmentFile={{ llm_router_env_file }}
ExecStart={{ llm_router_venv }}/bin/litellm --config {{ llm_router_config_file }} --host 0.0.0.0 --port {{ llm_router_api_port }}
# Restart policy is managed by the systemd_restart_policy role (99-restart-policy.conf).

[Install]
WantedBy=multi-user.target
