# roles/open_webui/templates/open-webui.service.j2
{{ ansible_managed | comment }}
[Unit]
Description=Open WebUI (LLM chat front-end)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ open_webui_user }}
Group={{ open_webui_user }}
WorkingDirectory={{ open_webui_home }}
EnvironmentFile=/etc/open-webui.env
ExecStart={{ open_webui_venv }}/bin/open-webui serve --host {{ open_webui_listen_host }} --port {{ open_webui_port }}
# Restart policy is managed by the systemd_restart_policy role (99-restart-policy.conf).

[Install]
WantedBy=multi-user.target
