# roles/traefik/templates/traefik.service.j2
{# Managed by Ansible (traefik role). #}
[Unit]
Description=Traefik reverse proxy / TLS ingress
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ traefik_user }}
Group={{ traefik_group }}
EnvironmentFile={{ traefik_acme_env_file }}
ExecStart={{ traefik_bin_path }} --configfile={{ traefik_config_dir }}/traefik.yml
Restart=on-failure
RestartSec=5
# Bind 80/443 without running as root.
AmbientCapabilities=CAP_NET_BIND_SERVICE
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
# Traefik must write the renewed cert to acme.json under the acme dir.
ReadWritePaths={{ traefik_acme_dir }}

[Install]
WantedBy=multi-user.target
