# roles/hermes_agent/templates/hermes-brain-watchdog.service.j2
{{ ansible_managed | comment }}
# One probe cycle of the Hermes brain-health watchdog. Fired by
# hermes-brain-watchdog.timer; runs as the hermes user so `hermes cron
# pause/resume` sees the same $HERMES_HOME the gateway uses.
[Unit]
Description=Hermes brain-health watchdog (pause/resume the cron fleet on a brain outage)
After=network-online.target hermes-gateway.service
Wants=network-online.target

[Service]
Type=oneshot
User={{ hermes_agent_user }}
Group={{ hermes_agent_user }}
Environment=HERMES_HOME={{ hermes_agent_hermes_home }}
# Same secrets file the gateway loads: HERMES_AGENT_MODEL_API_KEY (probe auth) +
# SLACK_BOT_TOKEN / SLACK_ALLOWED_USERS (DM alert). Leading '-' keeps the probe
# runnable while .env is still empty.
EnvironmentFile=-{{ hermes_agent_hermes_home }}/.env
ExecStart={{ hermes_agent_brain_watchdog_script_path }}
