# roles/hermes_agent/templates/hermes-brain-watchdog.timer.j2
{{ ansible_managed | comment }}
# Drive the Hermes brain-health watchdog every
# {{ hermes_agent_brain_watchdog_interval }} (mirrors the service_deadman cadence).
[Unit]
Description=Periodic Hermes brain-health probe

[Timer]
OnBootSec={{ hermes_agent_brain_watchdog_interval }}
OnUnitActiveSec={{ hermes_agent_brain_watchdog_interval }}
AccuracySec=5s

[Install]
WantedBy=timers.target
