# roles/hermes_agent/templates/hermes-gateway.service.j2
{{ ansible_managed | comment }}
[Unit]
Description=Hermes Agent gateway (autonomous agent daemon + cron + Kanban dispatcher)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ hermes_agent_user }}
Group={{ hermes_agent_user }}
WorkingDirectory={{ hermes_agent_home }}
Environment=HERMES_HOME={{ hermes_agent_hermes_home }}
# Load the deployed secrets so config.yaml's ${VAR} interpolation resolves. The
# leading '-' keeps the daemon startable while .env is still empty (no messaging
# platform wired yet) — the LiteLLM router master key still arrives via this file.
EnvironmentFile=-{{ hermes_agent_hermes_home }}/.env
ExecStart={{ hermes_agent_gateway_cmd }}
# Restart policy is managed by the systemd_restart_policy role (99-restart-policy.conf).

[Install]
WantedBy=multi-user.target
