# roles/technitium_install/templates/dns.service.j2
{{ ansible_managed | comment }}
# Technitium DNS Server — bare-metal .NET install. The app distribution is
# mirrored from MinIO (technitium_install role); the runtime is the ASP.NET Core
# runtime from the Microsoft apt repo. Layout matches the existing primary.
[Unit]
Description=Technitium DNS Server
After=network.target

[Service]
WorkingDirectory={{ technitium_install_app_dir }}/dns
ExecStart=/usr/bin/dotnet {{ technitium_install_app_dir }}/dns/DnsServerApp.dll {{ technitium_install_data_dir }}
Restart=always
# Restart service after 10 seconds if the dotnet service crashes:
RestartSec=10
SyslogIdentifier=dns-server

[Install]
WantedBy=multi-user.target
