# roles/media_monitor/templates/plex-session-probe.service.j2
[Unit]
Description=Plex Direct-Play session probe (alert on forced transcode)
Documentation=https://github.com/JacobPEvans/ansible-proxmox-apps
After=network-online.target

[Service]
Type=oneshot
ExecStart={{ media_monitor_bin_dir }}/plex-session-probe.sh
