# roles/media_monitor/templates/plex-codec-audit.service.j2
[Unit]
Description=Plex weekly library codec audit (flag non-AppleTV audio)
Documentation=https://github.com/JacobPEvans/ansible-proxmox-apps
After=network-online.target

[Service]
Type=oneshot
# Read-only ffprobe sweep — nice + idle IO so it never competes with live
# playback on a shared node.
Nice=15
IOSchedulingClass=idle
ExecStart={{ media_monitor_bin_dir }}/plex-codec-audit.sh
