# roles/media_monitor/templates/plex-codec-audit.timer.j2
[Unit]
Description=Weekly Plex library codec audit

[Timer]
OnCalendar={{ media_monitor_audit_oncalendar }}
AccuracySec=1h
Persistent=true

[Install]
WantedBy=timers.target
